# php5.6-apache/supervisor/supervisord.conf
[program:cron]
command = cron -f

[program:apache2]
command = apache2ctl -D FOREGROUND
